# x_engine/com.jryszardnoszczyk.linkedin-pull-user.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.jryszardnoszczyk.linkedin-pull-user</string>

    <!-- Weekly Sunday 07:00 — per-creator content updates slowly, weekly
         is cost-controlled. Per master plan v13 §3.4 + §7.3. Cost-cap
         fires per-call (`--max-cu 200`); weekly total bounded by user count. -->
    <key>StartCalendarInterval</key>
    <dict>
        <key>Weekday</key>
        <integer>0</integer>
        <key>Hour</key>
        <integer>7</integer>
        <key>Minute</key>
        <integer>0</integer>
    </dict>

    <key>WorkingDirectory</key>
    <string>/Users/jryszardnoszczyk/Documents/GitHub/gofreddy</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/Users/jryszardnoszczyk/Documents/GitHub/gofreddy/.venv/bin:/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin</string>
    </dict>

    <key>ProgramArguments</key>
    <array>
        <string>/bin/zsh</string>
        <string>-c</string>
        <string>./x_engine/run-linkedin-user.sh</string>
    </array>

    <key>StandardOutPath</key>
    <string>/Users/jryszardnoszczyk/Documents/GitHub/gofreddy/x_engine/logs/linkedin-user.log</string>

    <key>StandardErrorPath</key>
    <string>/Users/jryszardnoszczyk/Documents/GitHub/gofreddy/x_engine/logs/linkedin-user.err</string>

    <key>RunAtLoad</key>
    <false/>
</dict>
</plist>
